# mkosi.extra/etc/systemd/system/ssh@.service.d/override.conf
# 3717 – Stricter sshd absolute path requirement breaks xinetd
# https://bugzilla.mindrot.org/show_bug.cgi?id=3717
[Service]
ExecStart=
ExecStart=/usr/sbin/sshd -i -o UsePAM=no
